# deploy/com.toby.ride-dispatch-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.toby.ride-dispatch-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>/Users/toby/.cloudflared/ride-dispatch.yml</string>
        <string>run</string>
        <string>ride-dispatch</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/toby/ride-dispatch/logs/tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/toby/ride-dispatch/logs/tunnel.err</string>
</dict>
</plist>
