# deploy/com.toby.ride-dispatch-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.toby.ride-dispatch-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/accountname/ride-dispatch/.venv/bin/python</string>
        <string>-m</string>
        <string>ride_dispatch.bot</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/accountname/ride-dispatch</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/accountname/ride-dispatch/logs/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/accountname/ride-dispatch/logs/bot.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
